# .config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/home/simon/.config/guix/profiles/emacs/emacs/bin/notmuch new
