# .config/systemd/user/smtp.timer
[Unit]
Description=Mailbox send timer

[Timer]
OnBootSec=10m
OnUnitActiveSec=10m
Unit=smtp.service

[Install]
WantedBy=timers.target
